# nukkitx.service
[Unit]
Description=nukkitx mc bedrock edition server
After=network-online.target

[Service]
User=replace
WorkingDirectory=dirname/nukkitx
Type=forking
ExecStart=/bin/bash dirname/nukkitx/start.sh
ExecStop=/bin/bash dirname/nukkitx/stop.sh
GuessMainPID=no
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
